# roles/nats-server/templates/nats.service
[Unit]
Description=NATS messaging server

[Service]
ExecStart=/srv/nats/bin/gnatsd
Restart=on-failure

[Install]
WantedBy=multi-user.target